# ndpi/debian/xt-ndpi-filter.ndpi-filter.service
[Unit]
Description=Apply nDPI iptables rules (xt_ndpi)
Documentation=file:///etc/default/ndpi-filter
# Wants (not Requires) docker: the default DOCKER-USER chain needs dockerd, but
# the chain is configurable, so don't hard-fail on hosts without Docker.
After=docker.service
Wants=docker.service
# Re-apply whenever Docker (re)starts, since dockerd recreates DOCKER-USER.
PartOf=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/sbin/modprobe xt_ndpi
ExecStart=/usr/sbin/ndpi-filter start
ExecStop=/usr/sbin/ndpi-filter stop

[Install]
WantedBy=multi-user.target
